# ese3005/workspace/workspace/udev/vish_mmc_drive_mount.service
[Unit]
Description=mount the memory card drive!

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/mmc_vish_mount.sh

[Install]
WantedBy=multi-user.target
